# etc/init.d/npd6
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/bin/npd6
}

stop() {
	killall npd6
}
