# etc/init.d/radvd
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/bin/radvd
}

stop() {
	killall radvd
}
